# ops/systemd/chalo-analytics-updater.timer
[Unit]
Description=Chalo Dashboard - Analytics docs updater (every 10 minutes)

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target

